# Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:easync]
command=/root/.local/bin/easync /opt/server/easync.cfg
